# conductor.service
[Unit]
Description=Conductor Kanban
After=network.target

[Service]
Type=simple
User=ryan
WorkingDirectory=/home/ryan/Code/conductor
ExecStart=/usr/bin/node server.js
Environment=NODE_ENV=production
Environment=PORT=6194
Restart=on-failure

[Install]
WantedBy=multi-user.target
