# dispatcher.service
[Unit]
Description=Dispatcher Kanban
After=network.target

[Service]
Type=simple
User=ryan
WorkingDirectory=/home/ryan/Code/dispatcher
ExecStart=/home/ryan/.nvm/versions/node/v24.2.0/bin/node server.js
Environment=NODE_ENV=production
Environment=PORT=6194
Restart=on-failure

[Install]
WantedBy=multi-user.target
