# dispatch.service
[Unit]
Description=Dispatch Kanban
After=network.target

[Service]
Type=simple
User=ryan
WorkingDirectory=/home/ryan/Code/dispatch
ExecStart=/usr/bin/node server.js
Environment=NODE_ENV=production
Environment=PORT=3000
Restart=on-failure

[Install]
WantedBy=multi-user.target
