# orchestrel.service
[Unit]
Description=Orchestrel Kanban
After=network.target opencode.service
Wants=opencode.service

[Service]
Type=simple
User=ryan
WorkingDirectory=/home/ryan/Code/orchestrel
ExecStart=/home/ryan/.nvm/versions/node/v24.2.0/bin/pnpm dev
Environment=NODE_ENV=development
Environment=CF_TEAM_DOMAIN=wednesday-access
Environment=PATH=/home/ryan/.nvm/versions/node/v24.2.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=PORT=6194
Restart=on-failure

[Install]
WantedBy=multi-user.target
